# deploy/miper-alert@.service
[Unit]
Description=miper failure alert for %i

[Service]
Type=oneshot
User=miper
ExecStart=/home/miper/miper/deploy/alert-failure.sh %i
